# etc/udev/rules.d/99-undock-suspend.rules
# Only trigger on display disconnection events
ACTION=="remove", SUBSYSTEM=="drm", RUN+="/usr/local/bin/dock-handler.sh"
ACTION=="add", SUBSYSTEM=="drm", RUN+="/usr/local/bin/dock-handler.sh"
